# templates/default/opennms-pris.service.erb
[Install]
WantedBy=multi-user.target

[Unit]
Description=OpenNMS Provisioning Integration Server

[Service]
WorkingDirectory=<%=@app_home%>
ExecStart=/usr/bin/java -cp <%=@app_home%>/lib/*:<%=@app_home%>/opennms-pris.jar org.opennms.pris.Starter
ExecStop=/bin/kill -TERM $MAINPID
User=root
